# systemd/headroom-minimax.service
[Unit]
Description=Headroom proxy #2 (porta 8790) verso MiniMax — PROTETTO, non killare mai
Documentation=file://$HOME/Obsidian/Memoria
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=$HOME
Environment=ANTHROPIC_TARGET_API_URL=https://api.minimax.io/anthropic
Environment=HEADROOM_TELEMETRY=off
Environment=HEADROOM_UPDATE_CHECK=off
Environment=HEADROOM_SKIP_UPSTREAM_CHECK=1
ExecStart=/opt/headroom/.venv/bin/python3 -m headroom.cli proxy --host 127.0.0.1 --port 8790 --mode token --backend anthropic
Restart=always
RestartSec=2
OOMScoreAdjust=-900
OOMPolicy=continue

[Install]
WantedBy=default.target
